# chainflip-monitor.service
[Unit]
Description=Chainflip Validator Monitor Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/chainflip-monitor-bot
ExecStart=node src/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 